# docker/admin/emulate.cron
*/1 * * * * /usr/bin/python3 /app/emulate.py
